# deploy/systemd/spotifypod.service
[Unit]
Description=sPot - iPod Spotify UI
# Needs the X server (started via startx/openbox) plus Redis and a Connect device.
After=graphical.target redis-server.service raspotify.service
Wants=redis-server.service

[Service]
Type=simple
User=pi
# Adjust this path to wherever you cloned the repo:
WorkingDirectory=/home/pi/ipod_spotify_rasp2w/frontend
# SPOTIPY_CLIENT_ID / SECRET / REDIRECT_URI live here (KEY=value, no `export`).
EnvironmentFile=/home/pi/ipod_spotify_rasp2w/frontend/.env
Environment=DISPLAY=:0
ExecStart=/usr/bin/python3 spotifypod.py
Restart=always
RestartSec=3

[Install]
WantedBy=graphical.target
